# digger/bot/com.digger.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent for digger.

  Install:
    cp com.digger.bot.plist ~/Library/LaunchAgents/
    # Edit the copy: set UserName, the python path, and the working dir.
    launchctl load -w ~/Library/LaunchAgents/com.digger.bot.plist

  Check status:
    launchctl list | grep digger
    tail -f ~/Library/Logs/digger/stdout.log

  Unload:
    launchctl unload -w ~/Library/LaunchAgents/com.digger.bot.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.digger.bot</string>

    <!-- Absolute paths — no ~ or $HOME expansion in launchd. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jc/Nest/digger/bot/.venv/bin/python</string>
        <string>-u</string>
        <string>/Users/jc/Nest/digger/bot/slack_app.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jc/Nest/digger/bot</string>

    <!-- launchd starts processes with a minimal PATH that doesn't include
         Homebrew. The pipeline shells out to `goplaces` (installed via
         Homebrew), so we must extend PATH here. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Start on login and keep it running. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Don't hammer restart if it crashes in a loop (min 10s between tries). -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Logs. Create ~/Library/Logs/digger/ before loading. -->
    <key>StandardOutPath</key>
    <string>/Users/jc/Library/Logs/digger/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jc/Library/Logs/digger/stderr.log</string>

    <!-- Keep the agent idle-running at low CPU priority; bumps back up when active. -->
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
